# hw2/ansible/roles/yarn/templates/resourcemanager.service.j2
[Unit]
Description=Hadoop YARN ResourceManager
After=network.target

[Service]
User={{ hadoop_user }}
Group={{ hadoop_group }}
Environment=HADOOP_CONF_DIR=/etc/hadoop
ExecStart=/bin/bash -lc '{{ hadoop_symlink }}/bin/yarn resourcemanager'
Restart=on-failure

[Install]
WantedBy=multi-user.target
